# templates/dashing.upstart.erb
# Ubuntu upstart file at /etc/init/dashing.conf
#
# Generated by puppet: DO NOT EDIT!
#
# Author: Ricardo Rocha <ricardo@catalyst.net.nz>
#
start on runlevel [2345]
 
stop on runlevel [06]
 
pre-start script
    ENABLED="true"    # default value
    [ -f /etc/default/dashing ] && . /etc/default/dashing
    if [ "$ENABLED" = "true" ]
    then
        for dir in /var/run
        do
            if [ ! -d "$dir/dashing" ]
            then
                mkdir -p "$dir/dashing"
                chown <%= scope.lookupvar('dashing::run_user') %>:<%= scope.lookupvar('dashing::run_group') %> "$dir/dashing"
            fi
        done
 
        for CONFIG in /etc/dashing.conf /etc/dashing.d/*.conf
        do
            if [ -f "$CONFIG" ]
            then
                start dashing-instance CONFIG="$CONFIG" || :
            fi
        done
    fi
end script
 
post-stop script
    for CONFIG in `initctl list|sed -n 's/^dashing-instance (\(.*\)).*$/\1/p'`
    do
        if [ -f "$CONFIG" ]
        then
            stop dashing-instance CONFIG="$CONFIG" || :
        fi
    done
end script
